# deploy/ableducacion-backend.service
[Unit]
Description=ABL Educacion Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ableducacion/backend
Environment=PATH=/home/ubuntu/ableducacion/backend/venv/bin:/usr/bin
ExecStart=/home/ubuntu/ableducacion/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
